# debian/webengine-team.postinst
#!/bin/sh
# PostInstall script for WebEngine.

set -e

#DEBHELPER#

if [ "$1" = "configure" ]; then
	webserver="apache2"

	# Reload the webserver.
	if [ -x /usr/sbin/invoke-rc.d ]; then
		invoke-rc.d $webserver reload 3> /dev/null || true
	else
		/etc/init.d/$webserver reload 3> /dev/null || true
	fi
fi

exit 0
